# deploy/systemd/sirhosp-census-orchestrator.service
[Unit]
Description=SIRHOSP — Adaptive census orchestrator (continuous loop)
Documentation=https://github.com/carlosapgomes/sirhosp
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sirhosp
# Inicia o serviço dedicado census_orchestrator (profile orchestrator).
# O serviço roda run_adaptive_census_cycles --loop dentro do container.
ExecStart=docker compose -f compose.yml -f compose.prod.yml --profile orchestrator up --abort-on-container-exit census_orchestrator

# Para o serviço dedicado sem afetar web/worker.
ExecStop=docker compose -f compose.yml -f compose.prod.yml --profile orchestrator stop census_orchestrator

# Stale recovery is enabled by default in --loop mode.
# To disable temporarily, add: --disable-stale-recovery
# See deploy/README.md section 7 for details.

# Restart automatically on failure or crash
Restart=on-failure
RestartSec=10

# Graceful shutdown: SIGTERM is caught by the command's signal handler
TimeoutStopSec=30

# Logs go to journal (accessible via journalctl -u sirhosp-census-orchestrator)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sirhosp-census-orchestrator

[Install]
WantedBy=multi-user.target
